# vastbase-support.service
[Unit]
Description=Vastbase Intelligent Database Management Platform
Documentation=https://github.com/vastbase-support
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/fanrui/VastAIDemo
EnvironmentFile=/home/fanrui/VastAIDemo/.env
ExecStart=/usr/bin/node /home/fanrui/VastAIDemo/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vastbase-support

LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
